# system/create_deb_package/new_app_deb_package/DEBIAN/prerm
#!/bin/bash
set -e

# Stop and disable the service before removing the package
systemctl stop new_app.service
systemctl disable new_app.service
